# recipes-core/base-files/base-files/etc/rc.local
# This file is for adding custom commands which would be executed right after the system init finished.
# Currently it is empty by default.

exit 0
